# src/jts-serf/DEBIAN/postinst
#!/bin/bash
set -x -e

pam-auth-update --package

systemctl daemon-reload || true
deb-systemd-helper enable serf.service || true
deb-systemd-invoke try-restart serf.service || true
